# supervisor-workeropendb.conf
[program:gearmand]
priority = 100
command = python2.7 /opendb/worker_allschemas.py
redirect_stderr = true
stdout_logfile = /opendb/output.log
startsecs = 5
autorestart = true
